# cron.d/flexget
# m h dom mon dow user  command
0 *   * * *   abc    /usr/local/bin/flexget execute -c /config/config.yml --loglevel debug --now 2>&1 | logger -t FlexGet
